# debian/wlanpi-bluetooth.bt-agent.service
[Unit]
Description=Bluetooth Auth Agent
Requires=bluetooth.target
After=bluetooth.target
BindsTo=bluetooth.service
 
[Service]
ExecStart=/bin/sh -c "bt-agent --capability=NoInputNoOutput"
ExecStartPost=/bin/sh -c "bt-adapter --set Alias `hostname`"
ExecStartPost=/bin/sh -c "bt-adapter --set Powered `cat /etc/wlanpi-bluetooth/state`"
ExecStartPost=/bin/sh -c "bt-adapter --set Pairable 0"
TimeoutStopSec=5
Type=simple
 
[Install]
WantedBy=multi-user.target
